# packaging/kochab-agent.service
[Unit]
Description=Kochab Agent - Infrastructure monitoring daemon
Documentation=https://kochab.ai/docs/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/kochab-agent
WatchdogSec=120
Restart=on-failure
RestartSec=5
# Exit 70 = node decommissioned by the platform (410 GONE on poll). systemd
# must NOT restart in that case — the agent already wrote /etc/kochab/decommissioned,
# the boot guard would re-emit 70 on every restart anyway, and a restart loop
# would spam audit logs. Story 2.6.1 AC-12 (d).
RestartPreventExitStatus=70

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/etc/kochab

# Resource limits
LimitNOFILE=65536
MemoryMax=128M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kochab-agent

[Install]
WantedBy=multi-user.target
