# packaging/kochab-agent.service
[Unit]
Description=Kochab Agent - Infrastructure monitoring daemon
Documentation=https://kochab.ai/docs/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/kochab-agent
WatchdogSec=120
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/etc/kochab

# Resource limits
LimitNOFILE=65536
MemoryMax=128M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kochab-agent

[Install]
WantedBy=multi-user.target
